# nsm-dummy/systemd/com.contiautomotive.NodeStateManager.service.in
[Service]
Type = notify
BusName = com.contiautomotive.NodeStateManager.service
ExecStart = @libdir@/nsm-dummy-@BOOT_MANAGER_VERSION_API@/nsm-dummy
