# nsm-dummy/systemd/nsm-dummy.service.in
[Service]
Type = dbus
BusName = com.contiautomotive.NodeStateManager
ExecStart = @libdir@/boot-manager-@BOOT_MANAGER_VERSION_API@/nsm-dummy
