# tests/legacy-app-handler/legacy-app-handler-test1.service.in
[Service]
ExecStart = @libdir@/boot-manager-@BOOT_MANAGER_VERSION_API@/legacy-app-handler-test Test1
ExecStartPost = @libdir@/boot-manager-@BOOT_MANAGER_VERSION_API@/legacy-app-handler \
                --register --unit=legacy-app-handler-test1.service
