# tests/legacy-app-handler/legacy-app-handler-test2.service.in
[Service]
ExecStart = @libdir@/node-startup-controller-@NODE_STARTUP_CONTROLLER_VERSION_API@/legacy-app-handler-test Test2
ExecStartPost = @libdir@/node-startup-controller-@NODE_STARTUP_CONTROLLER_VERSION_API@/legacy-app-handler \
                --register --unit=legacy-app-handler-test2.service
