# node-startup-controller/node-startup-controller.service.in
[Service]
Type = notify
ExecStart = @libdir@/node-startup-controller-@NODE_STARTUP_CONTROLLER_VERSION_API@/node-startup-controller
WatchdogSec = 5
