# legacy-app-handler/systemd/org.genivi.LegacyAppHandler1.service.in
[Service]
BusName = org.genivi.LegacyAppHandler1.service
ExecStart = @libdir@/legacy-app-handler-@BOOT_MANAGER_VERSION_API@/legacy-app-handler
